# analyzer.timer
[Unit]
Description=Run analyzer.py every half an hour

[Timer]
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target