# analyzer.timer
[Unit]
Description=Run analyzer.py every half an hour

[Timer]
OnCalendar=*:0/30
Persistent=True

[Install]
WantedBy=timers.target